# setup/adhan.service
[Unit]
Description=Adhan Prayer Time Playback Service
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/adhan
ExecStart=/opt/adhan/venv/bin/python -m src.main
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PULSE_RUNTIME_PATH=/run/user/1000/pulse

# Hardening
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
